# systemd/system/redis.service.d/local.conf
[Service]
RuntimeDirectory=redis
RuntimeDirectoryMode=750
