# systemd/system/mastodon-streaming.service.d/local.conf
[Service]
RuntimeDirectory=mastodon-streaming
RuntimeDirectoryMode=750
Environment="SOCKET=/run/mastodon-streaming/mastodon-streaming.sock"
