# systemd/system/mastodon-sidekiq.service.d/override.conf
[Unit]
After=syslog-ng@default.service
Requires=syslog-ng@default.service

[Service]
StandardOutput=file:/run/mastodon-sidekiq-log
StandardError=file:/run/mastodon-sidekiq-log
