# systemd/system/mastodon-sidekiq.service.d/override.conf
[Service]
StandardOutput=file:/run/mastodon-sidekiq-log
StandardError=file:/run/mastodon-sidekiq-log
