# systemd/system/valkey.service.d/local.conf
[Service]
RuntimeDirectory=valkey
RuntimeDirectoryMode=750
